# scripts/launchd/com.fuqing.add-orders-index.one-shot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fuqing.add-orders-index.one-shot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/hutou/homebrew/bin/python3</string>
        <string>/Users/hutou/Desktop/fuqin-date/fuqing-crm-analytics/scripts/etl/add_orders_pay_time_user_id_index.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/fuqing-add-orders-index.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/fuqing-add-orders-index.err</string>
</dict>
</plist>
