# scripts/launchd/com.fuqing.memory-size-monitor.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fuqing.memory-size-monitor.weekly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/hutou/homebrew/bin/python3</string>
        <string>/Users/hutou/Desktop/fuqin-date/fuqing-crm-analytics/scripts/ops/memory_size_monitor.py</string>
    </array>

    <!-- Sprint 201+ R7: 每周日 04:15 跑 1 次 (跟 R6 错开 15 分钟) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>15</integer>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/fuqing-memory-size.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/fuqing-memory-size.err</string>
</dict>
</plist>