# scripts/launchd/com.fuqing.pre-existing-fail-monitor.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fuqing.pre-existing-fail-monitor.weekly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/hutou/homebrew/bin/python3</string>
        <string>/Users/hutou/Desktop/fuqin-date/fuqing-crm-analytics/scripts/pre_existing_fail_monitor.py</string>
    </array>

    <!-- Sprint 201+ R6: 每周日 04:00 跑 1 次 (跨 sprint 0 debt stable 监控) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/fuqing-pre-existing-fail.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/fuqing-pre-existing-fail.err</string>
</dict>
</plist>