# scripts/launchd/com.fuqing.duckdb-to-parquet-etl.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fuqing.duckdb-to-parquet-etl.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/hutou/homebrew/bin/python3</string>
        <string>/Users/hutou/Desktop/fuqin-date/fuqing-crm-analytics/scripts/etl/duckdb_to_parquet_etl.py</string>
        <string>--output-dir</string>
        <string>/Users/hutou/Desktop/fuqin-date/fuqing-crm-analytics/data/processed/postgresql16_parquet</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>40</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/fuqing-duckdb-to-parquet-etl.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/fuqing-duckdb-to-parquet-etl.err</string>
</dict>
</plist>
