# scripts/launchd/com.fuqing.db-size-alert.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fuqing.db-size-alert.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/hutou/homebrew/bin/python3</string>
        <string>/Users/hutou/Desktop/fuqin-date/fuqing-crm-analytics/scripts/ops/check_db_size.py</string>
    </array>

    <!-- Sprint 201 R2 L2.6: 每天 04:00 跑 1 次 (跟 duckdb-backup.daily 03:30 错开) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/fuqing-db-size-alert.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/fuqing-db-size-alert.err</string>
</dict>
</plist>
